# deploy/t4-shadow/cihash.service
[Unit]
Description=CIHash T4 shadow GitHub App webhook service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f cihash
ExecStart=/usr/bin/docker run --name cihash --network agent-webhook-hub_web --user 997:987 --read-only --cap-drop=ALL --security-opt=no-new-privileges --pids-limit=256 --memory=256m --cpus=1 --tmpfs /tmp:rw,noexec,nosuid,size=16m --env-file /srv/cihash/secrets/cihash.env --mount type=bind,src=/srv/cihash/config,dst=/srv/cihash/config,readonly --mount type=bind,src=/srv/cihash/secrets,dst=/srv/cihash/secrets,readonly --mount type=bind,src=/srv/cihash/receipts,dst=/srv/cihash/receipts --mount type=bind,src=/srv/cihash/state,dst=/srv/cihash/state sha256:1d1e81a86ae541a6223cb41721bb5e0aec81734dab686a1519ac3b3e8f0384b3 serve --config /srv/cihash/config/hosted.json
ExecStop=/usr/bin/docker stop --time 10 cihash
ExecStopPost=-/usr/bin/docker rm -f cihash
Restart=on-failure
RestartSec=5s
TimeoutStartSec=0
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
